# matrix/rootfs/etc/services.d/matrix/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Matrix
# Runs the Matrix Daemon
# ==============================================================================

bashio::log.level "$(bashio::config 'log_level' 'warning')"

config=/config/addons_config/matrix

bashio::log.info "Stop Matrix Service"
/opt/venvs/matrix-synapse/bin/synctl stop "$config"/homeserver.yaml

if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "Matrix crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "Service restart after closing"
