# matrix/rootfs/etc/services.d/matrix/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Matrix
# Runs the Matrix Daemon
# ==============================================================================

bashio::log.level "$(bashio::config 'log_level' 'warning')"

config=/config/addons_config/matrix

bashio::log.info "Starting Matrix Service..."
if [ -f /share/matrix/homeserver.pid ]; then
    rm /share/matrix/homeserver.pid || true
fi

#export PYTHONPATH=/opt/venvs/matrix-synapse/lib/python3.9/site-packages/
#exec python3 -B -m synapse.app.homeserver -c "$config"/homeserver.yaml
/opt/venvs/matrix-synapse/bin/synctl start "$config"/homeserver.yaml
