# matrix/rootfs/etc/services.d/matrix/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Matrix
# Runs the Matrix Daemon
# ==============================================================================

bashio::log.level "$(bashio::config 'log_level' 'warning')"

config=/config/addons_config/matrix

bashio::log.info "Starting Matrix Service..."
rm /share/matrix/homeserver.pid || true
exec /opt/venvs/matrix-synapse/bin/synctl start "$config"/homeserver.yaml
