# dyndns/rootfs/etc/services.d/dyndns/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: DynDns
# Runs the DynDns
# ==============================================================================

bashio::log.level "$(bashio::config 'log_level' 'warning')"

# Load Options and Start DynDns
bashio::log.info "Load DynDns Configuration..."
if bashio::config.has_value 'api_token'; then

    if bashio::config.has_value 'domains'; then

        bashio::log.trace "Read API Token"
        apiToken=$(bashio::config 'api_token')
        export HETZNER_AUTH_API_TOKEN="${apiToken}"

        bashio::log.trace "Read Interval"
        interval=$(bashio::config 'interval' '1h')

        while true; do
            for domainAsJson in $(bashio::config 'domains'); do
                domain=$(bashio::jq "${domainAsJson}" ".domain")
                bashio::log.info "Execute DynDns Update for Domain '${domain}'"

                for recordAsJson in $(bashio::jq "${domainAsJson}" ".records[]"); do
                    record=$(bashio::jq "${recordAsJson}" ".record")
                    ttl=$(bashio::jq "${recordAsJson}" ".ttl")
                    type=$(bashio::jq "${recordAsJson}" ".type")

                    bashio::log.info "Update Record '${record}' with TTL ${ttl} and Type ${type}"
                    /usr/bin/bashio /etc/dyndns/dyndns.sh -Z "${domain}" -n "${record}" -t "${ttl}" -T "${type}"
                done
            done

            sleep "${interval}"
        done
    fi

else
    bashio:log.fatal "No API Token found."
    bashio::exit.nok
fi
