# dyndns/rootfs/etc/services.d/dyndns/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: DynDns
# Runs the DynDns
# ==============================================================================

bashio::log.level "$(bashio::config 'log_level' 'warning')"

# Load Options and Start DynDns
bashio::log.info "Load DynDns Configuration..."
if bashio::config.has_value 'api_token'; then

    if bashio::config.has_value 'records'; then

        bashio::log.trace "Read API Token"
        apiToken=$(bashio::config 'api_token')
        export HETZNER_AUTH_API_TOKEN="${apiToken}"

        bashio::log.trace "Read Interval"
        interval=$(bashio::config 'interval' '1h')

        while true; do
            for recordAsJson in $(bashio::config 'records'); do
                domain=$(bashio::jq "${recordAsJson}" ".domain")
                record=$(bashio::jq "${recordAsJson}" ".name")
                ttl=$(bashio::jq "${recordAsJson}" ".ttl")
                type=$(bashio::jq "${recordAsJson}" ".type")

                bashio::log.info "Execute DynDns Update for Domain '${domain}'"
                bashio::log.info "Update Record '${record}' with TTL ${ttl} and Type ${type}"

                /usr/bin/bashio /etc/dyndns/dyndns.sh -Z "${domain}" -n "${record}" -t "${ttl}" -T "${type}"
            done

            sleep "${interval}"
        done
    fi

else
    bashio::log.fatal "No API Token found."
    bashio::exit.nok
fi
