# mongodb/rootfs/etc/services.d/mongodb/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: DynDns
# Runs the DynDns
# ==============================================================================

# Set Log Level
bashio::log.level "$(bashio::config 'log_level' 'warning')"

MONGODB_DATA=/data/databases
MONGODB_LOG=/data/logs

echo -n "${MONGODB_DATA}" > /var/run/s6/container_environment/MONGODB_DATA
echo -n "${MONGODB_LOG}" > /var/run/s6/container_environment/MONGODB_LOG

# Create Datafolder
mkdir -p "${MONGODB_DATA}"
mkdir -p "${MONGODB_LOG}"

if [ -f "${MONGODB_DATA}/mongodb.err" ]; then
    rm -f "${MONGODB_DATA}/mongodb.err"
fi

# Redirect errors
ln -s /proc/1/fd/1 "${MONGODB_DATA}/mongodb.err"

bashio::log.info "Start MongoDB ..."
exec mongod --dbpath "${MONGODB_DATA}" --logpath "${MONGODB_LOG}"/mongod.log
