# guacamole/rootfs/etc/services.d/guacd/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Guacamole
# Runs the Guacamole Daemon interface
# ==============================================================================

prog="guacd"
exec="/usr/local/sbin/$prog"
pidfile="/var/run/$prog.pid"

bashio::log.level "$(bashio::config 'log_level' 'warning')"

bashio::log.info "Starting the Guacamole Daemon..."
exec "$exec" -p "$pidfile" -f
