# dns/rootfs/etc/services.d/dns/run
#!/command/with-contenv bashio
# shellcheck disable=SC2155
# ==============================================================================
# Home Assistant Community Add-on: Dns/Dhcp Server
# Runs the Dns/Dhcp Server
# ==============================================================================

# Set Log Level
if bashio::config.has_value 'log_level'; then
    bashio::log.info 'Set Log Level'
    log_level=$(bashio::config 'log_level')
    bashio::log.level "${log_level}"
fi

# Set some Environment Variables if available
bashio::log.info "Load Dns/Dhcp Server Configuration..."
bashio::config.has_value 'domain' && export DNS_SERVER_DOMAIN=$(bashio::config 'domain') && bashio::log.info 'Domain set'
bashio::config.has_value 'admin_password' && export DNS_SERVER_ADMIN_PASSWORD=$(bashio::config 'admin_password') && bashio::log.info 'Admin Password set'
bashio::config.has_value 'prefer_ipv6' && export DNS_SERVER_PREFER_IPV6=$(bashio::config 'prefer_ipv6') && bashio::log.info 'Prefer IPv6 set'
bashio::config.has_value 'dns_over_http' && export DNS_SERVER_OPTIONAL_PROTOCOL_DNS_OVER_HTTP=$(bashio::config 'dns_over_http') && bashio::log.info 'Dns over Http set'
bashio::config.has_value 'recursion' && export DNS_SERVER_RECURSION=$(bashio::config 'recursion') && bashio::log.info 'Recursion set'
bashio::config.has_value 'recursion_denied_networks' && export DNS_SERVER_RECURSION_DENIED_NETWORKS=$(bashio::config 'recursion_denied_networks') && bashio::log.info 'Recursion denied Networks set'
bashio::config.has_value 'recursion_allowed_networks' && export DNS_SERVER_RECURSION_ALLOWED_NETWORKS=$(bashio::config 'recursion_allowed_networks') && bashio::log.info 'Recursion allowed Networks set'
bashio::config.has_value 'enable_blocking' && export DNS_SERVER_ENABLE_BLOCKING=$(bashio::config 'enable_blocking') && bashio::log.info 'Blocking set'
bashio::config.has_value 'allow_txt_blocking_report' && export DNS_SERVER_ALLOW_TXT_BLOCKING_REPORT=$(bashio::config 'allow_txt_blocking_report') && bashio::log.info 'Text Blocking set'
bashio::config.has_value 'block_list_urls' && export DNS_SERVER_BLOCK_LIST_URLS=$(bashio::config 'block_list_urls') && bashio::log.info 'Blocking List set'
bashio::config.has_value 'forwarders' && export DNS_SERVER_FORWARDERS=$(bashio::config 'forwarders') && bashio::log.info 'Forwarders set'
bashio::config.has_value 'forwarder_protocol' && export DNS_SERVER_FORWARDER_PROTOCOL=$(bashio::config 'forwarder_protocol') && bashio::log.info 'Forwareder Protocol set'
bashio::config.has_value 'local_time' && export DNS_SERVER_LOG_USING_LOCAL_TIME=$(bashio::config 'local_time') && bashio::log.info 'Local Time set'

currentIPv4=$(hostname -i)
currentIPv6=$(hostname -I | grep -oP '[a-f0-9:]+:+[a-f0-9]+')

bashio::log.info "Starting Dns/Dhcp Server..."
bashio::log.info "Use IP '${currentIPv4}' for Home Assistant IPv4 DNS Server Konfiguration"
bashio::log.info "Use IP '${currentIPv6}' for Home Assistant IPv6 DNS Server Konfiguration"
exec /usr/bin/dotnet /opt/dns/DnsServerApp.dll /data/config
