# vocechat/rootfs/etc/s6-overlay/s6-rc.d/vochechat-core/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash

# Set Log Level
bashio::log.level "$(bashio::config 'log_level' 'warning')"

bashio::log.info "Start VoceChat ..."
exec vocechat-server /etc/vocechat/config.toml --daemon

