# vocechat/rootfs/etc/s6-overlay/s6-rc.d/vocechat-pre/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
ROOT=/data
DATA_ROOT=/data/db
NEW_INSTALL=false

if ! bashio::fs.directory_exists "${DATA_ROOT}"; then
    NEW_INSTALL=true
fi

bashio::config.has_value 'domain' && DOMAIN=$(bashio::config 'domain') && bashio::log.info 'Domain set'

# Save variables
echo -n "${DATA_ROOT}" > /var/run/s6/container_environment/DATA_ROOT
echo -n "${NEW_INSTALL}" > /var/run/s6/container_environment/NEW_INSTALL
echo -n "${DOMAIN}" > /var/run/s6/container_environment/DOMAIN

# Create folders
mkdir -p "${DATA_ROOT}"

# Redirect log output
rm -f "${ROOT}/app.err"
ln -s /proc/1/fd/1 "${ROOT}/app.err"
