# technitium/rootfs/etc/services.d/technitium/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs Technitium
# ==============================================================================

bashio::log.info "Starting Technitium DNS Server ..."
exec /etc/dns/start.sh
