# guacamole/rootfs/etc/services.d/guacamole/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Guacamole
# Runs the Guacamole interface
# ==============================================================================

bashio::log.level "$(bashio::config 'log_level' 'warning')"

bashio::log.info "Starting the Guacamole Interface..."
exec /opt/tomcat/bin/startup.sh
